# chef/cookbooks/pfda/templates/pfda-gsrs.service.erb
[Unit]
Description=pfda-gsrs

[Service]
EnvironmentFile=<%= @gsrs_environment_file %>
ExecStart=<%= @gsrs_run_script %>
ExecStop="/bin/kill `cat <%= @gsrs_path %>/RUNNING_PID`"
WorkingDirectory=<%= @gsrs_path %>

StandardOutput=file:/srv/www/precision_fda/current/log/ginas.log
StandardError=file:/srv/www/precision_fda/current/log/ginas.log

SyslogIdentifier=ginas

Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
