# chef/cookbooks/pfda/templates/sidekiq.service.erb
[Unit]
Description=<%= @name %>

[Service]
Environment=<%= @environment.map { |key, val| %Q{"#{key}=#{val}"} }.join(' ') %>
ExecStartPre=/bin/bash -c 'rm /srv/www/precision_fda/current/log/sidekiq.log'
ExecStart=<%= @command %>
ExecReload=/bin/kill -<%= @reload_signal %> $MAINPID
KillSignal=<%= @stop_signal %>
User=<%= @user %>
WorkingDirectory=<%= @directory %>
Restart=<%= @restart_mode %>

StandardOutput=file:/srv/www/precision_fda/current/log/sidekiq.log
StandardError=file:/srv/www/precision_fda/current/log/sidekiq.log

SyslogIdentifier=sidekiq

[Install]
WantedBy=multi-user.target
